# source/Config/ScheduledTask/Disable QPCore Service.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>ZSkycat</Author>
    <URI>\Disable QPCore Service</URI>
  </RegistrationInfo>
  <Triggers>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription>&lt;QueryList&gt;&lt;Query Id="0" Path="Application"&gt;&lt;Select Path="Application"&gt;*[System[Provider[@Name='QPCore'] and EventID=0]]&lt;/Select&gt;&lt;/Query&gt;&lt;/QueryList&gt;</Subscription>
    </EventTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-21-0000000000-0000000000-0000000000-0000</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>false</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT1H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-WindowStyle Hidden -EncodedCommand ZgB1AG4AYwB0AGkAbwBuACAAQwBsAGUAYQByAEcAYQByAGIAYQBnAGUAIAB7AA0ACgAgACAAIAAgAHQAYQBzAGsAawBpAGwAbAAuAGUAeABlACAALwBmACAALwB0ACAALwBpAG0AIABRAFEAUAByAG8AdABlAGMAdAAuAGUAeABlAA0ACgAgACAAIAAgAFMAZQB0AC0AUwBlAHIAdgBpAGMAZQAgAFEAUABDAG8AcgBlACAALQBTAHQAYQByAHQAdQBwAFQAeQBwAGUAIABNAGEAbgB1AGEAbAANAAoAfQANAAoADQAKAFsAZABvAHUAYgBsAGUAXQAkAG0AYQB4AFcAcwAgAD0AIAAxADAAMAANAAoAWwBkAGEAdABlAHQAaQBtAGUAXQAkAGIAZQBnAGkAbgBUAGkAbQBlAFcAcwAgAD0AIABHAGUAdAAtAEQAYQB0AGUADQAKAFsAZABhAHQAZQB0AGkAbQBlAF0AJABiAGUAZwBpAG4AVABpAG0AZQBOAHUAbABsACAAPQAgAEcAZQB0AC0ARABhAHQAZQANAAoAUwB0AGEAcgB0AC0AUAByAG8AYwBlAHMAcwAgACcAQwA6AFwAUAByAG8AZwByAGEAbQAgAEYAaQBsAGUAcwAgACgAeAA4ADYAKQBcAFQAZQBuAGMAZQBuAHQAXABUAEkATQBcAEIAaQBuAFwAUQBRAFMAYwBMAGEAdQBuAGMAaABlAHIALgBlAHgAZQAnAA0ACgANAAoAdwBoAGkAbABlACAAKAAkAHQAcgB1AGUAKQAgAHsADQAKACAAIAAgACAAVwByAGkAdABlAC0ATwB1AHQAcAB1AHQAIAAnAFMAbABlAGUAcAAnAA0ACgAgACAAIAAgAFMAdABhAHIAdAAtAFMAbABlAGUAcAAgADEADQAKAA0ACgAgACAAIAAgACQAbgBvAHcAIAA9ACAARwBlAHQALQBEAGEAdABlAA0ACgAgACAAIAAgACQAcAByAG8AYwBlAHMAcwAgAD0AIABHAGUAdAAtAFAAcgBvAGMAZQBzAHMAIAAnAHQAaQBtACcADQAKAA0ACgAgACAAIAAgAGkAZgAgACgAJABwAHIAbwBjAGUAcwBzACkAIAB7AA0ACgAgACAAIAAgACAAIAAgACAAJABiAGUAZwBpAG4AVABpAG0AZQBOAHUAbABsACAAPQAgACQAbgBvAHcADQAKAA0ACgAgACAAIAAgACAAIAAgACAAIwAgAHt2VV+MWxBiwGhLbToAIABTXyAAVABpAG0AIACFUVhb3U8BYzN6mlv2ZQ0ACgAgACAAIAAgACAAIAAgACAAJABjAHUAcgByAGUAbgB0AE0AYQB4ACAAPQAgACgAJABwAHIAbwBjAGUAcwBzACAAfAAgAE0AZQBhAHMAdQByAGUALQBPAGIAagBlAGMAdAAgAC0AUAByAG8AcABlAHIAdAB5ACAAVwBTACAALQBNAGEAeABpAG0AdQBtACkALgBNAGEAeABpAG0AdQBtACAALwAgADEAbQBiAA0ACgANAAoAIAAgACAAIAAgACAAIAAgAGkAZgAgACgAWwBTAHkAcwB0AGUAbQAuAE0AYQB0AGgAXQA6ADoAQQBiAHMAKAAkAGMAdQByAHIAZQBuAHQATQBhAHgAIAAtACAAJABtAGEAeABXAHMAKQAgAC0AbAB0ACAANQApACAAewANAAoAIAAgACAAIAAgACAAIAAgACAAIAAgACAAaQBmACAAKAAoACQAbgBvAHcAIAAtACAAJABiAGUAZwBpAG4AVABpAG0AZQBXAHMAKQAuAFMAZQBjAG8AbgBkAHMAIAAtAGcAdAAgADUAKQAgAHsADQAKACAAIAAgACAAIAAgACAAIAAgACAAIAAgACAAIAAgACAAQwBsAGUAYQByAEcAYQByAGIAYQBnAGUADQAKACAAIAAgACAAIAAgACAAIAAgACAAIAAgACAAIAAgACAAVwByAGkAdABlAC0ATwB1AHQAcAB1AHQAIAAnAEMAbABlAGEAcgAgACgAdwBzACkAJwANAAoAIAAgACAAIAAgACAAIAAgACAAIAAgACAAIAAgACAAIABiAHIAZQBhAGsADQAKACAAIAAgACAAIAAgACAAIAAgACAAIAAgAH0ADQAKACAAIAAgACAAIAAgACAAIAB9AA0ACgAgACAAIAAgACAAIAAgACAAZQBsAHMAZQAgAHsADQAKACAAIAAgACAAIAAgACAAIAAgACAAIAAgACQAYgBlAGcAaQBuAFQAaQBtAGUAVwBzACAAPQAgACQAbgBvAHcADQAKACAAIAAgACAAIAAgACAAIAB9AA0ACgANAAoAIAAgACAAIAAgACAAIAAgAGkAZgAgACgAJABjAHUAcgByAGUAbgB0AE0AYQB4ACAALQBnAHQAIAAkAG0AYQB4AFcAcwApACAAewAgACQAbQBhAHgAVwBzACAAPQAgACQAYwB1AHIAcgBlAG4AdABNAGEAeAAgAH0ADQAKACAAIAAgACAAfQANAAoAIAAgACAAIABlAGwAcwBlACAAewANAAoAIAAgACAAIAAgACAAIAAgACMAIAACXzhewGhLbToAIABTXyAAVABpAG0AIADbjwt6DU5YWyhX9mUNAAoAIAAgACAAIAAgACAAIAAgAGkAZgAgACgAKAAkAG4AbwB3ACAALQAgACQAYgBlAGcAaQBuAFQAaQBtAGUATgB1AGwAbAApAC4AUwBlAGMAbwBuAGQAcwAgAC0AZwB0ACAANQApACAAewANAAoAIAAgACAAIAAgACAAIAAgACAAIAAgACAAQwBsAGUAYQByAEcAYQByAGIAYQBnAGUADQAKACAAIAAgACAAIAAgACAAIAAgACAAIAAgAFcAcgBpAHQAZQAtAE8AdQB0AHAAdQB0ACAAJwBDAGwAZQBhAHIAIAAoAG4AdQBsAGwAKQAnAA0ACgAgACAAIAAgACAAIAAgACAAIAAgACAAIABiAHIAZQBhAGsADQAKACAAIAAgACAAIAAgACAAIAB9AA0ACgAgACAAIAAgAH0ADQAKAH0ADQAKAA==</Arguments>
    </Exec>
  </Actions>
</Task>